# systemd/user/compton.service
[Unit]
Description=compositor

[Service]
ExecStart=/usr/bin/compton --backend glx
Restart=always

