# systemd/user/picom.service
[Unit]
Description=compositor

[Service]
ExecStart=/usr/bin/picom --vsync --backend xrender
Restart=always

