# systemd/user/udiskie.service
[Unit]
Description=automounter

[Service]
ExecStart=/usr/bin/udiskie -a -T --no-appindicator -F -0 -N
Restart=always

