# systemd/user/urxvtd.service
[Service]
Descrption=urxvtd

[Unit]
ExecStart=/usr/bin/urxvtd
Restart=always

