# systemd/user/urxvtd.service
[Unit]
Descrption=urxvtd

[Service]
ExecStart=/usr/bin/urxvtd
Restart=always

