# systemd/user/xbindkeys.service
[Unit]
Description=map trackball buttons

[Service]
ExecStart=/usr/bin/xbindkeys -n
Restart=always

