# systemd/user/xcompmgr.service
[Unit]
Description=compositor

[Service]
ExecStart=/usr/bin/xcompmgr -c
Restart=always

